# wireguard/wg-sync.path
[Unit]
Description=Watch for WireGuard peer config changes

[Path]
PathChanged=/opt/wassal/wireguard/peers/.trigger
Unit=wg-sync.service

[Install]
WantedBy=multi-user.target
