# wireguard/wg-sync.service
[Unit]
Description=Sync WireGuard peers from Wassal backend

[Service]
Type=oneshot
ExecStart=/opt/wassal/wireguard/sync-peers.sh
